# deploy/agent-character-watchdog.service
[Unit]
Description=ACK Monitor Watchdog (root-owned)
After=agent-character-monitor.service agent-enforcer.service
Requires=agent-character-monitor.service agent-enforcer.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lib/agent-character-kit/ack_watchdog.py
Restart=always
RestartSec=3
User=root
StandardOutput=journal
StandardError=journal
Environment=ENFORCER_SOCKET=/run/agent-enforcer/main.sock
Environment=AGENT_WORKSPACE=/var/lib/agent-character-kit/workspace
Environment=ACK_MONITOR_PID=/var/lib/agent-character-kit/ack-monitor.pid
Environment=ACK_WATCHDOG_PID=/var/lib/agent-character-kit/ack-watchdog.pid
Environment=ACK_MONITOR_BIN=/usr/local/lib/agent-character-kit/ack_monitor.py

[Install]
WantedBy=multi-user.target
