# deploy/agent-character-monitor.service
[Unit]
Description=ACK Acknowledgment Monitor (root-owned)
After=agent-enforcer.service
Requires=agent-enforcer.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lib/agent-character-kit/ack_monitor.py
Restart=always
RestartSec=3
User=root
StandardOutput=journal
StandardError=journal
Environment=ENFORCER_SOCKET=${ENFORCER_SOCKET:-/run/agent-enforcer/main.sock}
Environment=AGENT_WORKSPACE=${AGENT_WORKSPACE:-/var/lib/agent-character-kit/workspace}
Environment=ACK_ACK_LOG=${ACK_ACK_LOG:-/tmp/agent-character-kit-ack.jsonl}
Environment=ACK_MONITOR_PID=${ACK_MONITOR_PID:-/var/lib/agent-character-kit/ack-monitor.pid}
Environment=ACK_MONITOR_STATE=${ACK_MONITOR_STATE:-/var/lib/agent-character-kit/ack-monitor.pos}

[Install]
WantedBy=multi-user.target
